# xiaozhi-addon/rootfs/etc/s6-overlay/s6-rc.d/svc-xiaozhi/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Xiaozhi ESP32 Server - S6 v3 longrun service
# ==============================================================================

declare -r SERVER_DIR="/opt/xiaozhi-esp32-server"
declare -r CONFIG_FILE="${SERVER_DIR}/data/.config.yaml"

# Read add-on options
ANTHROPIC_API_KEY=$(bashio::config 'anthropic_api_key')
MODEL_NAME=$(bashio::config 'model_name')
LANGUAGE=$(bashio::config 'language')
TTS_VOICE=$(bashio::config 'tts_voice')
PROMPT=$(bashio::config 'prompt')

# Get HA Supervisor token (automatically provided by HAOS)
HA_TOKEN="${SUPERVISOR_TOKEN}"

# Determine TTS voice based on language if not explicitly set
if [ -z "${TTS_VOICE}" ]; then
    case "${LANGUAGE}" in
        fr) TTS_VOICE="fr-FR-HenriNeural" ;;
        en) TTS_VOICE="en-US-GuyNeural" ;;
        zh) TTS_VOICE="zh-CN-XiaoxiaoNeural" ;;
        *)  TTS_VOICE="fr-FR-HenriNeural" ;;
    esac
fi

# Build HA devices string from options (semicolon-separated: "location,name,entity_id;...")
HA_DEVICES=""
if bashio::config.has_value 'ha_devices'; then
    for device in $(bashio::config 'ha_devices | .[]'); do
        if [ -n "${HA_DEVICES}" ]; then
            HA_DEVICES="${HA_DEVICES};${device}"
        else
            HA_DEVICES="${device}"
        fi
    done
fi

if [ -z "${HA_DEVICES}" ]; then
    HA_DEVICES="salon,example_device,light.example"
fi

# Get local IP
LOCAL_IP=$(hostname -I | awk '{print $1}')
if [ -z "${LOCAL_IP}" ]; then
    LOCAL_IP="0.0.0.0"
fi

# Indent prompt
PROMPT_INDENTED=$(echo "${PROMPT}" | sed 's/^/  /')

bashio::log.info "Generating server configuration..."
bashio::log.info "  LLM Model: ${MODEL_NAME}"
bashio::log.info "  Language: ${LANGUAGE}"
bashio::log.info "  TTS: Local Piper (fr_FR-siwis-medium)"
bashio::log.info "  Local IP: ${LOCAL_IP}"

mkdir -p "${SERVER_DIR}/data"

cat > "${CONFIG_FILE}" << ENDOFCONFIG
# Auto-generated by Xiaozhi HAOS add-on
server:
  ip: 0.0.0.0
  port: 8000
  http_port: 8003
  websocket: "ws://${LOCAL_IP}:8000/xiaozhi/v1/"
  vision_explain: "http://${LOCAL_IP}:8003/mcp/vision/explain"
  timezone_offset: +1
  auth:
    enabled: true
    allowed_devices:
      - "3c:0f:02:de:a8:18"
  auth_key: "$(python3 -c 'import secrets; print(secrets.token_hex(32))')"
  mqtt_gateway: null
  mqtt_signature_key: null
  udp_gateway: null

selected_module:
  VAD: SileroVAD
  ASR: SherpaASR
  LLM: ClaudeLLM
  TTS: LocalTTS
  Memory: nomem
  Intent: function_call

LLM:
  ClaudeLLM:
    type: Claude
    model_name: "${MODEL_NAME}"
    api_key: "${ANTHROPIC_API_KEY}"
    max_tokens: 4096
    temperature: 0.7

TTS:
  LocalTTS:
    type: sherpa_tts
    model_dir: models/vits-piper-fr_FR-siwis-medium
    model_file: fr_FR-siwis-medium.onnx
    speed: 1.0
    speaker_id: 0
    output_dir: tmp/

ASR:
  SherpaASR:
    type: sherpa_onnx_local
    model_dir: models/sherpa-onnx-nemo-parakeet-tdt-0.6b-v3-int8
    output_dir: tmp/
    model_type: transducer

VAD:
  SileroVAD:
    type: silero
    threshold: 0.5
    threshold_low: 0.3
    model_dir: models/snakers4_silero-vad
    min_silence_duration_ms: 200

Intent:
  function_call:
    type: function_call
    functions:
      - change_role
      - hass_get_state
      - hass_set_state
      - hass_play_music

plugins:
  home_assistant:
    devices: "${HA_DEVICES}"
    base_url: "http://supervisor/core"
    api_key: "${HA_TOKEN}"

prompt: |
${PROMPT_INDENTED}
ENDOFCONFIG

bashio::log.info "Starting Xiaozhi ESP32 Server..."

cd "${SERVER_DIR}"
exec python3 app.py
